# build/supervisord.conf
[supervisord]
nodaemon=false

[program:elasticsearch]
command=/usr/local/bin/docker-entrypoint.sh eswrapper
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true

[program:exceptionless]
command=dotnet Exceptionless.Web.dll
directory=/app
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
